# deploy/systemd/gpu-ready.service
[Unit]
Description=Wait for NVIDIA GPU to be ready after boot
Documentation=https://github.com/ggml-org/llama.cpp
DefaultDependencies=no
Before=docker.service

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/opt/llama/scripts/gpu-ready.sh
TimeoutStartSec=180
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
