# deploy/systemd/llama@.service
[Unit]
Description=llama.cpp server (%I)
Documentation=https://github.com/ggml-org/llama.cpp
# gpu-ready.service delays start until NVIDIA GPU is fully initialized,
# preventing degraded throughput ("Post-Reboot Throughput Incident").
Requires=docker.service gpu-ready.service
After=docker.service network-online.target gpu-ready.service
StartLimitIntervalSec=120
StartLimitBurst=3

[Service]
# llama.sh starts a detached Docker container and exits immediately.
# RemainAfterExit=yes keeps systemd tracking the service as "active".
Type=oneshot
RemainAfterExit=yes
ExecStart=/opt/llama/llama.sh start %i
ExecStop=/opt/llama/llama.sh stop
TimeoutStopSec=30

# Security hardening
ProtectSystem=strict
ReadOnlyPaths=/opt/llama
ReadWritePaths=/opt/llama/models
BindPaths=/var/run/docker.sock:/var/run/docker.sock
PrivateTmp=true
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
