# deploy/systemd/llama-postboot-restart@.timer
[Unit]
Description=Timer for post-boot restart of llama.cpp server (%I)
Documentation=https://github.com/ggml-org/llama.cpp
# Run once after boot, 2 min after llama@%i starts
After=llama@%i.service

[Timer]
OnBootSec=2min
Unit=llama-postboot-restart@%i.service

[Install]
WantedBy=multi-user.target
