# deploy/systemd/llama@.service
[Unit]
Description=llama.cpp server (%I)
Documentation=https://github.com/ggml-org/llama.cpp
# nvidia-persistenced.service keeps GPU driver state initialized across
# reboots, preventing the "Post-Reboot Throughput Incident" (degraded ~1.5
# tok/s on first model load). gpu-ready.service is an additional safety check.
Requires=docker.service nvidia-persistenced.service gpu-ready.service
After=docker.service network-online.target nvidia-persistenced.service gpu-ready.service
StartLimitIntervalSec=120
StartLimitBurst=3

[Service]
# llama.sh starts a detached Docker container and exits immediately.
# RemainAfterExit=yes keeps systemd tracking the service as "active".
Type=oneshot
RemainAfterExit=yes
ExecStart=/opt/llama/llama.sh start %i
ExecStop=/opt/llama/llama.sh stop
TimeoutStopSec=30

# Security hardening
ProtectSystem=strict
ReadOnlyPaths=/opt/llama
ReadWritePaths=/opt/llama/models
BindPaths=/var/run/docker.sock:/var/run/docker.sock
PrivateTmp=true
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
