# deploy/systemd/llama-gpu-watchdog.timer
[Unit]
Description=Timer for llama.cpp GPU watchdog

[Timer]
OnBootSec=2min
OnUnitActiveSec=2min

[Install]
WantedBy=timers.target
